# com.organization.psm-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    LaunchDaemon: com.organization.psm-monitor
    Description : Runs psm_monitor.zsh hourly to collect FileVault unlock-attempt
                  and backoff state for all local users (UID > 500) into a SQLite3 database.

    Installation:
        sudo cp com.organization.psm-monitor.plist /Library/LaunchDaemons/
        sudo chmod 644 /Library/LaunchDaemons/com.organization.psm-monitor.plist
        sudo chown root:wheel /Library/LaunchDaemons/com.organization.psm-monitor.plist
        sudo cp psm_monitor.zsh /usr/local/bin/psm_monitor.zsh
        sudo chmod 755 /usr/local/bin/psm_monitor.zsh
        sudo chown root:wheel /usr/local/bin/psm_monitor.zsh
        sudo launchctl load /Library/LaunchDaemons/com.organization.psm-monitor.plist

    Removal:
        sudo launchctl unload /Library/LaunchDaemons/com.organization.psm-monitor.plist
        sudo rm /Library/LaunchDaemons/com.organization.psm-monitor.plist
        sudo rm /usr/local/bin/psm_monitor.zsh

    Customize:
        - Replace "com.organization" with your organization's reverse-DNS identifier.
        - Adjust StartCalendarInterval to change the collection schedule.
        - DB is written to /var/db/psm_monitor/psm_monitor.db
-->
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.organization.psm-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/usr/local/bin/psm_monitor.zsh</string>
    </array>

    <!-- Run at the top of every hour -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run immediately when the daemon is loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Runs as root (required for psm list) — no UserName key needed for LaunchDaemon -->

    <key>StandardOutPath</key>
    <string>/var/log/psm_monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/psm_monitor.log</string>

</dict>
</plist>
