# pkgsource/postinstall
#!/bin/bash

DAEMON_LABEL="gov.nasa.cset.psm-monitor"
PLIST="/Library/LaunchDaemons/${DAEMON_LABEL}.plist"
SCRIPT="/usr/local/bin/psm_monitor.zsh"
DB_DIR="/var/db/psm_monitor"

# Ensure /usr/local/bin exists (may not on a fresh system)
mkdir -p /usr/local/bin

# Set ownership and permissions on installed files
chown root:wheel "$SCRIPT"
chmod 755 "$SCRIPT"

chown root:wheel "$PLIST"
chmod 644 "$PLIST"

# Create the database directory with appropriate permissions
mkdir -p "$DB_DIR"
chown root:wheel "$DB_DIR"
chmod 755 "$DB_DIR"

# Load the LaunchDaemon
echo "Loading ${DAEMON_LABEL}..."
launchctl load "$PLIST"

exit 0
